# supervisord.conf
[supervisord]
nodaemon=true

[program:bro]
command=/opt/bro/bin/broctl deploy

[eventlistener:exit_on_any_fatal]
command=exit-event-listener
events=PROCESS_STATE_FATAL
